# .config/systemd/user/antigravity-sync.service
[Unit]
Description=Antigravity IDE Conversation Sync Daemon
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/sidsrlal/.gemini/antigravity-ide/scratch/antigravity-sync-daemon.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
